# bender-worker.service
[Unit]
Description=bender-worker
User=bender

[Service]
Type=simple
ExecStart=/usr/bin/bender-worker

[Install]
WantedBy=multi-user.target
